# systemd/system/dnf-automatic-install.timer.d/time.conf
[Timer]
OnBootSec=
OnCalendar=tue 18:00
